# systemd/smct-scanner.service
[Unit]
Description=SMCT BTCUSDT Scanner Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/workspace/AlertBot
ExecStart=/usr/bin/python3 /workspace/AlertBot/scripts/run_scanner.py
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
